# deploy/systemd/mumble-bg-control.service
[Unit]
Description=mumble-bg HTTP control server
After=network.target postgresql.service

[Service]
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__ENV_FILE__
Environment=BG_ENV_FILE=__ENV_FILE__
Environment=BG_KEY_DIR=__KEY_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__VENV_DIR__/bin/python -I -m bg.control_main --noreload
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
